# templates/pulpcore-worker@.service.erb
[Unit]
Description=Pulp Worker
After=network-online.target
Wants=network-online.target postgresql.service

[Service]
Type=simple
Environment="DJANGO_SETTINGS_MODULE=pulpcore.app.settings"
Environment="PULP_SETTINGS=<%= scope['pulpcore::settings_file'] %>"
User=<%= scope['pulpcore::user'] %>
Group=<%= scope['pulpcore::group'] %>
WorkingDirectory=<%= scope['pulpcore::user_home'] %>
RuntimeDirectory=pulpcore-worker-%i
ExecStart=/usr/bin/pulpcore-worker
SyslogIdentifier=pulpcore-worker-%i

# This provides reconnect support for PostgreSQL and Redis. Without reconnect support, if either
# is not available at startup or becomes disconnected, this process will die and not respawn.
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

